# root/etc/services.d/laravel-echo-server/run
#!/usr/bin/with-contenv bash

cd /app && laravel-echo-server start